# service/symfonychat.target
[Unit]
Description="Symfonychat consumers"
Requires=symfonychat.consumer@webpush.service

[Install]
WantedBy=multi-user.target